# config/systemd/waybar-auto-hide.service.d/override.conf
[Service]
ExecStart=
ExecStart=/home/keita/.config/hypr/scripts/waybar_auto_hide --always-hidden
